# project_code/wifipi_router/etc/supervisor/conf.d/router.conf
[program:router]
command = python /var/www/router/web.pyc #--port=80
autostart = true
startsecs = 5
user = root
redirect_stderr = true
